# meta-genivi-demo/recipes-extended/browser-poc/browser-poc-git/browser.service
[Unit]
Description=Genivi Browser PoC service
Requires=dbus.service

[Service]
ExecStart=/opt/browser/bin/browser
