# meta-genivi-dev/meta-genivi-dev/recipes-core/systemd/systemd/focused.target
[Unit]
Description=Responsible for starting GENIVI components directly being used for the user context.
After=basic.target
