# meta-genivi-dev/meta-genivi-dev/recipes-core/systemd/systemd/lazy.target
[Unit]
Description=Responsible for starting lazy background applications.
After=unfocused.target
