# meta-genivi-dev/meta-genivi-dev/recipes-core/systemd/systemd/unfocused.target
[Unit]
Description=Contains all applications in the system.
After=focused.target
