# recipes-demo-hmi/genivi-demo-platform-hmi/gdp-hmi-background/gdp-hmi-background.service
[Unit]
Description=GENIVI Demo Platform HMI - Background

[Service]
Environment=LD_PRELOAD=/usr/lib/libEGL.so
ExecStart=/usr/bin/gdp-hmi-background