# recipes-demo-hmi/genivi-demo-platform-hmi/gdp-hmi-launcher2/StartLauncher.service
[Unit]
Description=Genivi Start Launcher service on boot
Requires=dbus.service
After=AudioManager.service pulseaudio.service

[Service]
ExecStart=/usr/bin/start_launcher.sh
