# meta-genivi-dev/recipes-extended/chromium/chromium-wayland/chromium.service
[Unit]
Description=Chromium browser

[Service]
ExecStart=/usr/bin/google-chrome
