# recipes-demo-hmi/genivi-demo-platform-hmi/qml-example/qml-example.service
[Unit]
Description=Simple QML (ivi-shell) Application

[Service]
Environment=LD_PRELOAD=/usr/lib/libEGL.so
ExecStart=/usr/bin/qml-example