# meta-genivi-dev/recipes-dev-hmi/genivi-dev-platform-hmi/gdp-hmi-launcher2/StartLauncher.path
[Unit]
Description=GENIVI Start Launcher when weston is ready
Requires=weston.service
After=weston.service

[Path]
PathExists=/var/run/user/0/wayland-0
