# meta-genivi-dev/meta-genivi-dev/recipes-extended/automotive-message-broker/automotive-message-broker/ambd.service
#
# ambd systemd service unit file
#

[Unit]
Description=Automotive Message Broker
After=syslog.target

[Service]
Type=dbus
BusName=org.automotive.message.broker
ExecStart=/usr/bin/ambd --config /etc/ambd/logreplayerconfig.in.json

[Install]
WantedBy=multi-user.target
