# meta-genivi-dev/meta-genivi-dev/recipes-extended/automotive-message-broker/automotive-message-broker/ambd.service
#
# ambd systemd service unit file
#

[Unit]
Description=Automotive Message Broker
After=syslog.target

[Service]
Type=dbus
User=ambd
BusName=org.automotive.message.broker
ExecStart=/usr/bin/ambd

[Install]
WantedBy=multi-user.target
