# recipes-graphics/wayland/wayland-ivi-extension/EGLWLMockNavigation.service
[Unit]
Description=EGL Wayland Mock Navigation

[Service]
ExecStart=/usr/bin/EGLWLMockNavigation
