# recipes-graphics/wayland/wayland-ivi-extension/EGLWLInputEventExample.service
[Unit]
Description=EGL Wayland Input Event Example

[Service]
ExecStart=/usr/bin/EGLWLInputEventExample
