# recipes-demo-hmi/genivi-demo-platform-hmi/gdp-hmi-panel/gdp-hmi-panel.service
[Unit]
Description=GENIVI Demo Platform HMI - Panel
After=gdp-hmi-controller.service

[Service]
Environment=LD_PRELOAD=/usr/lib/libEGL.so
ExecStart=/usr/bin/gdp-hmi-panel
